# buildspecs/api-deploy/api-scripts/financial-news.service
[Unit]
Description=Financial News API Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/financial-news-engine
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/financial-news-engine/app.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/financial-news-engine/logs/service-output.log
StandardError=append:/opt/financial-news-engine/logs/service-error.log
SyslogIdentifier=financial-news

[Install]
WantedBy=multi-user.target
