# buildspecs/api-deploy/scripts/financial-news.service
[Unit]
Description=Financial News Engine Flask App
After=network.target

[Service]
User=root
WorkingDirectory=/opt/financial-news-engine
Environment="CORS_ALLOWED_ORIGINS=https://financialnewsengine.com,https://www.financialnewsengine.com"
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/financial-news-engine/.env
ExecStart=/usr/local/bin/gunicorn --workers 3 --bind 0.0.0.0:5000 app:app --access-logfile logs/access.log --error-logfile logs/error.log --capture-output --log-level info
Restart=always
TimeoutStartSec=60
RestartSec=10

[Install]
WantedBy=multi-user.target 