# buildspecs/api-deploy/scripts/financial-news.service
[Unit]
Description=Financial News API Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/financial-news
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/financial-news/venv/bin/python3 /opt/financial-news/app.py
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=financial-news

[Install]
WantedBy=multi-user.target 